# etc/services.d/main_run/run
#!/command/with-contenv bash
python3 /ext/main_script.py $TOKEN $ZONEID $DNSID $DNSNAME $PRIMARYURL $BOTTOKEN $CHATID